# walllust-daemon.service
[Unit]
Description=Walllust Wallpaper Daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=%h/.local/bin/walllust-daemon
Restart=always
RestartSec=3
Slice=session.slice

[Install]
WantedBy=default.target
